# systemd/simplelogin-handler.service
[Unit]
Description=SimpleLogin.io clone email handler daemon
After=network.target

[Service]
# User and group to run the service as
User=simplelogin
Group=simplelogin

# The working directory where the script is located
WorkingDirectory=/home/simplelogin/sl-app

ExecStart=/home/simplelogin/sl-app/.venv/bin/python3 email_handler.py

# Restart the service if it fails
Restart=on-failure
# Wait 5 seconds before attempting to restart
RestartSec=5

# Optional: Graceful reload on 'systemctl reload simplelogin-gunicorn'
# This sends a HUP signal to Gunicorn, which reloads workers without downtime.
ExecReload=/bin/kill -s HUP $MAINPID

# Standard output and error will be sent to the systemd journal
StandardOutput=journal
StandardError=journal

# Recommended security settings:
# Ensures the service gets its own /tmp and /var/tmp directories
PrivateTmp=true
# Prevents the service from gaining new privileges
NoNewPrivileges=true
# Limits the maximum number of open file descriptors, good for servers
LimitNOFILE=65536

[Install]
# This unit should be enabled when the system reaches multi-user runlevel
WantedBy=multi-user.target
